# roles/fileshare/templates/auto_iso@.service.j2
# {{ ansible_managed }}

[Unit]
Description=Auto ISO Service (%f)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/auto_iso --udev '%f'
